# templates/restic-backup.service.d/prepare.conf.j2
# {{ ansible_managed }}

[Unit]
Requires=restic-prepare.service
After=restic-prepare.service
